# systemd/dtlogin.service
[Unit]
Description=Start CDE Login

[Service]
Type=oneshot
ExecStart=/bin/sh -c /usr/dt/bin/dtlogin -daemon

[Install]
WantedBy=multi-user.target
